# backup/gitea-backup.service
[Unit]
RequiresMountsFor=/docker/home/gitea

[Service]
ExecStart=/usr/local/bin/gitea-backup.sh

