# backup/gitea-backup.timer
[Timer]
OnCalendar=20:55

[Install]
WantedBy=timers.target
